# infra/temporal/temporal-worker.service
[Unit]
Description=Temporal Worker
After=network.target temporal-server.service
Requires=temporal-server.service

[Service]
Type=simple
WorkingDirectory=/home/ian/instant
ExecStart=/home/ian/instant/temporal-worker
Restart=always
RestartSec=5
User=ian
Group=ian

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=temporal-worker

[Install]
WantedBy=multi-user.target

